# config/sh/profile.d/z.sh
test -f ~/.local/bin/z.sh && source ~/.local/bin/z.sh
