# config/zsh/profile.d/options.zsh
setopt autocd
setopt extendedglob
setopt nomatch
setopt completealiases
setopt autonamedirs
setopt listtypes
setopt nobeep
setopt nolistbeep
setopt nonotify
setopt nocorrect
setopt promptsubst
